# cat-detector.service
[Unit]
Description=Simple Cat Detector Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cat-detector
Environment=PATH=/home/pi/cat-detector/venv/bin
ExecStart=/home/pi/cat-detector/venv/bin/python /home/pi/cat-detector/start_detection.py
Restart=always
RestartSec=10
KillMode=mixed
RestartPreventExitStatus=1
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 